# scripts/launchd/co.tamazia.engine-cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>co.tamazia.engine-cycle</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION/scripts/run-engine-cycle.sh</string>
  </array>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION/logs/launchd.out</string>
  <key>StandardErrorPath</key><string>/Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION/logs/launchd.err</string>
</dict>
</plist>
