# scripts/W14-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.tamazia.w14-ad-intel</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd /Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION &amp;&amp; source .env &amp;&amp; node src/lib/ad-intel/W14-cron.js 30 &gt;&gt; reports/w14-ad-intel.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION/reports/w14-launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/amanigga/Desktop/TAMAZIA-REBUILD/COWORK-OS-EXECUTION/reports/w14-launchd.err</string>
</dict></plist>
